# data/systemd/deskbar.service
[Unit]
Description=Vitruvian Deskbar
Conflicts=getty@tty1.service
After=janus.service input_server.service notification_server.service
Requires=janus.service

[Service]
Type=forking
PIDFile=/run/vitruvian/Deskbar.pid
ExecStart=/system/servers/janus_launch Deskbar
StandardOutput=append:/var/log/deskbar.log
StandardError=append:/var/log/deskbar.log
Restart=no

[Install]
WantedBy=graphical.target
